# mempalace-viz.service
[Unit]
Description=MemPalace 3D Visualization Server
After=network.target
Wants=network.target

[Service]
Type=forking
WorkingDirectory=/home/iggut/.openclaw/workspace/mempalace-viz
ExecStartPre=-/usr/bin/pkill -f "mempalace-viz" || true
ExecStart=/usr/bin/node server.js &>/var/log/mempalace-viz.log & sleep 2
Restart=always
RestartSec=5
User=iggut
Group=iggut
Environment="NODE_ENV=production"

[Install]
WantedBy=default.target
